# spec/sample/crontab
#This is my Cron Comment

CRON_VARIABLE="a specific value"

* * * * * root /usr/local/bin/some_script
